# ansible/roles/linux/files/etc/cron.d/jmorris.julia
15 12 * * 1-5 jmorris DISPLAY=:0 google-chrome-stable --new-window https://excella.unanet.biz/excella/action/time
15 12 * * 1-5 jmorris DISPLAY=:0 google-chrome-stable --new-window https://outlook.office.com/mail/inbox
15 12 * * 1-5 jmorris DISPLAY=:0 google-chrome-stable --new-window https://calendar.google.com
15 12 * * 1-5 jmorris DISPLAY=:0 google-chrome-stable --new-window https://vocabulary.com/vocabtrainer
15 12 * * 1-5 jmorris DISPLAY=:0 if ps -A | grep -q slack; then wmctrl -r slack -b add,demands_attention; else slack; fi &
# make sure messaging windows are always up
0 * * * * jmorris DISPLAY=:0 /home/jmorris/dotfiles/scripts/launch_messaging_windows.sh
# to delete clipboard history
*/10 * * * * jmorris DISPLAY=:0 /home/jmorris/dotfiles/scripts/clear_clipboard.sh
# to backup cron 
0 13 * * * jmorris cp /etc/cron.d/jmorris /home/jmorris/dotfiles/ansible/roles/linux/files/etc/cron.d/jmorris.julia
